# service/sync/run
#!/bin/bash

nohup qrsboxcli sync >/dev/null 2>&1 &
exec qrsboxcli log 2>&1
